# mac-optimizer/com.user.mactmon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mactmon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mac/Projects/ThisPC/mac-optimizer/mac-telegram-monitor.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mac/Projects/ThisPC/mac-optimizer/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mac/Projects/ThisPC/mac-optimizer/monitor.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/mac/Projects/ThisPC/mac-optimizer</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MACTMON_BOT_TOKEN</key>
        <string>ISI_BOT_TOKEN_DISINI</string>
        <key>MACTMON_CHAT_ID</key>
        <string>ISI_CHAT_ID_DISINI</string>
    </dict>
</dict>
</plist>
